# src/build/clangd_init_d
#!/bin/sh

### BEGIN INIT INFO
# Provides:        clangd
# Required-Start:  $network $local_fs $syslog
# Required-Stop:   $network $local_fs $syslog
# Default-Start:   2 3 4 5
# Default-Stop:    1
# Short-Description: Start Clangd daemon
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/bin/dist-clang
NAME=clangd
DAEMON=/usr/bin/dist-clang/${NAME}
DESC="Distributed Clang Daemon"
DAEMON_ARGS="--config /etc/${NAME}.conf"
PIDFILE="/var/run/${NAME}.pid"

set -e

. /lib/lsb/init-functions

test -x $DAEMON || exit 0

clangd_start() {
  log_daemon_msg "Starting clangd"
  start-stop-daemon --start -p $PIDFILE -x $DAEMON -n $NAME --startas ${DAEMON}_wrapper -- $DAEMON_ARGS 
  log_end_msg $?
}

clangd_stop() {
  log_daemon_msg "Stopping clangd"
  start-stop-daemon -p $PIDFILE --stop --retry 3 || echo -n "...which is not running"
  log_end_msg $?
}

case "$1" in
  start|stop)
    ${NAME}_${1}
    ;;
  restart|reload|force-reload)
    if [ -s $PIDFILE ] && kill -0 $(cat $PIDFILE) >/dev/null 2>&1; then
      ${NAME}_stop
      ${NAME}_start
    fi
    ;;
  force-stop)
    ${NAME}_stop
    killall ${NAME} || true
    sleep 3
    killall -9 ${NAME} || true
    ;;
  status)
    status_of_proc -p $PIDFILE "$DAEMON" "system-wide ${NAME}" && exit 0 || exit $?
    ;;
  *)
    echo "Usage: /etc/init.d/${NAME} {start|stop|force-stop|restart|reload|force-reload|status}"
    exit 1
    ;;
esac

exit 0
